# deploy/sumariza-ai.service
[Unit]
Description=Sumariza AI - Twitter Content Scraper
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/sumariza-ai
ExecStart=/var/www/sumariza-ai/bin/sumariza
Restart=always
RestartSec=5

# Environment
EnvironmentFile=/var/www/sumariza-ai/.env

# Chrome/Chromium path - auto-detected during setup, stored in .env
# If not in .env, common paths are tried by chromedp

# Security hardening (relaxed for Chromium headless)
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/var/www/sumariza-ai /tmp
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sumariza-ai

[Install]
WantedBy=multi-user.target

